# bike-monitor.service
[Unit]
Description=Bike Theft Prevention System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi5
WorkingDirectory=/home/pi5/Desktop/BikeMonitor
ExecStart=/usr/bin/python3 /home/pi5/Desktop/BikeMonitor/hailo_theft_prevention.py
Restart=always
RestartSec=10
Environment="DISPLAY=:0"
Environment="PYTHONUNBUFFERED=1"
Environment="QT_QPA_PLATFORM=offscreen"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
